# systemd/bitcoin-monitor.service
[Unit]
Description=Bitcoin Monitor

[Service]
User=<USER>
ExecStart=/home/<USER>/start.sh

[Install]
WantedBy=multi-user.target
